# rootfs/etc/services.d/cron/run
#!/usr/bin/with-contenv bash

echo Starting cron...

crond -f \
    -c /etc/crontabs/ \
    -l ${CRON_LOG_LEVEL} \
    -L "${CRON_LOG_FILE}"
