# docker/rootfs/etc/services.d/cron/run
#!/usr/bin/with-contenv bash

echo Starting cron...
exec /usr/sbin/crond -f -S -l 5 -c /etc/crontabs
